# etc/systemd/gtk-theme.timer
[Unit]
Description=Change the GTK theme based on the current time.

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
